# ansible/playbooks/roles/x2_autoware/files/autoware.service
[Unit]
Description=autoware
After=network-online.target can-adm.service

[Service]
ExecStart=/opt/autoware/services/autoware/run.sh
KillSignal=SIGINT
User=autoware

[Install]
WantedBy=multi-user.target
