# ansible/playbooks/files/simple_system_launch/autoware_launch.service
[Unit]
Description=autoware
After=network-online.target autoware_api_launch.service autoware_localization_launch.service autoware_system_launch.service

[Service]
ExecStart=/opt/autoware/services/system_launch/autoware_launch.service/run.sh
KillSignal=SIGINT
User=autoware

[Install]
WantedBy=multi-user.target
